# Exercises/out/production/Exercises/Linux Service/rot13.service
[Unit]
Description=ROT13 demo service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=dhivyaruban
ExecStart=/usr/bin/env php /Desktop/server.php

[Install]
WantedBy=multi-user.target
